# etc/services.d/apache/run
#!/usr/bin/execlineb -P

exec apache2ctl -DFOREGROUND
